# supervisor/supervisord.conf
[supervisord]
nodaemon=true ; Ejecutar en primer plano para que el contenedor no se cierre

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;" ; Iniciar Nginx en primer plano
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:php-fpm]
command=php-fpm -F ; Iniciar PHP-FPM en primer plano
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0